# cfg/rc.local
#!/bin/bash
iptables-restore < /etc/iptables_tor
exit 0
